# setup/puppet/modules/gunicorn/files/app.conf
description "Server instance"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec python /opt/app/product/manage.py run_gunicorn >> /tmp/app.log 2>&1
